# 02_rootfs/rootfs/etc/inittab


::sysinit:/etc/init.d/rcS
console::askfirst:-/bin/sh
::ctrlaltdel:/sbin/reboot
::shoutdown:/bin/umount -a -r

